# rc.local
#!/bin/bash
#sysctl -w net.ipv4.ip_forward=1
iptables -t nat -N OFFLINEWEB
iptables -t nat -F OFFLINEWEB
#iptables -t nat -A OFFLINEWEB -p tcp ! -s 172.20.0.100 --dport 80 -j REDIRECT --to-port 3129
#iptables -t nat -A OFFLINEWEB -p tcp ! -s 172.20.0.100 --dport 443 -j REDIRECT --to-port 3130
iptables -t nat -D PREROUTING -j OFFLINEWEB
iptables -t nat -A PREROUTING -j OFFLINEWEB
iptables -t nat -D POSTROUTING -j MASQUERADE
#iptables -t nat -A POSTROUTING -j MASQUERADE
